# rootfs/etc/inittab
#etc/inittab
::sysinit:/etc/init.d/rcs

::console::askfirst:-/bin/sh
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a
